# tools/com.user.audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.<user>.audit</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string><vault>/tools/run-audit.sh</string>
    <string>auto</string>
  </array>

  <!-- Dimanche uniquement (Weekday 0), fenêtre 09:00→10:30, toutes les 15 min. Les gardes
       du wrapper (déjà-fait cette semaine, offline, OOO) font qu'il ne tourne qu'une fois. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>15</integer></dict>
    <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>15</integer></dict>
    <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string><vault>/tools/logs/audit.out.log</string>
  <key>StandardErrorPath</key>
  <string><vault>/tools/logs/audit.err.log</string>
</dict>
</plist>
